# templates/hapax_legomenon.service
[Unit]
Description=Foo ${unique_word}

[Service]
ExecStart=/usr/sbin/foo-daemon

[Install]
WantedBy=multi-user.target
